# supervisord.conf
; =============================================================================
; AETHERIS — supervisord configuration
; Gestiona FastAPI (interno :8000) y Streamlit (publico :7860) en un solo
; contenedor Docker para Hugging Face Spaces.
;
; Nota: supervisord corre como el usuario aetheris (uid=1000) ya que el
; Dockerfile establece USER aetheris antes de ENTRYPOINT. NO se incluye
; la directiva user= en [supervisord] para evitar errores de setuid().
; =============================================================================

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; ---------------------------------------------------------------------------
; FastAPI — servidor interno (no expuesto directamente)
; ---------------------------------------------------------------------------
[program:api]
command=python -m uvicorn aetheris.api.main:app --host 0.0.0.0 --port 8000 --workers 1 --log-level info --no-access-log
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/tmp/api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=100

; ---------------------------------------------------------------------------
; Streamlit — puerto publico 7860 (HF Spaces)
; startsecs=15 da tiempo al lifespan de FastAPI para arrancar completamente
; ---------------------------------------------------------------------------
[program:ui]
command=python -m streamlit run aetheris/ui/app.py --server.port=7860 --server.address=0.0.0.0 --server.headless=true --server.enableCORS=false --server.enableXsrfProtection=false --browser.gatherUsageStats=false
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=5
startsecs=15
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/tmp/ui.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=200
